# packer/files/puma.service
[Unit]
Description=Puma

[Service]
WorkingDirectory=/home/appuser/reddit
ExecStart=/usr/local/bin/puma

[Install]
WantedBy=multi-user.target
